# debian/sreview-detect.cron.d
*/30 * * * * sreview [ -x /usr/bin/sreview-detect ] && /usr/bin/sreview-detect
